# src-debian-package/DEBIAN/prerm
#!/bin/sh

systemctl stop deye2mqtt.service
systemctl disable deye2mqtt.service

exit 0
